# conf/supervisor-app.conf
[program:nginx-app]
command = /usr/sbin/nginx

[program:memcached]
command=/usr/bin/memcached -m 64 -p 11211 -u memcache -l 127.0.0.1 -DFOREGROUND

[program:python-app]
command= /src/server.sh
